# crontab
*/1 * * * * /usr/bin/git -C /nnix.com pull origin main
*/1 * * * * /usr/bin/git -C /nnix.com/capsule pull origin main