# crontab
*/1 * * * * cd /nnix.com && /usr/bin/git pull origin main
*/1 * * * * cd /nnix.com/capsule && /usr/bin/git pull origin main