# crontab
*/1 * * * * ec2-user /usr/bin/git --git-dir=/web/nnix.com pull >/dev/null 2>&1
*/1 * * * * ec2-user /usr/bin/git --git-dir=/web/nnix.com/capsule pull >/dev/null 2>&1